# .ddev/crontab
* * * * * php /var/www/html/bin/wp cron event run --due-now
